# App/munkipkg/Remote Access/scripts/postinstall
#!/bin/bash

launchctl load /Library/LaunchDaemons/com.github.remoteaccess.install.plist

exit 0
